# systemd/legale-bot.service
[Unit]
Description=Legale Bot Telegram Webhook Daemon
Documentation=https://github.com/legale/rag-telegram-gpt-bot
After=network.target

[Service]
Type=forking
User=legale-bot
Group=legale-bot
WorkingDirectory=/home/legale-bot/legale-bot

# Environment
Environment="PATH=/home/legale-bot/.local/bin:/usr/local/bin:/usr/bin:/bin"

# PID file
PIDFile=/home/legale-bot/legale-bot/profiles/default/bot.pid

# Command to run (using legale.py orchestrator)
ExecStart=/home/legale-bot/legale-bot/legale.py bot daemon --port 8000
ExecStop=/bin/kill -TERM $MAINPID

# Kill all child processes when stopping
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/legale-bot/legale-bot

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=legale-bot

[Install]
WantedBy=multi-user.target
